# deploy/opencarrier.service
[Unit]
Description=OpenCarrier Agent OS Daemon
Documentation=https://github.com/yinnho/opencarrier
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=opencarrier
Group=opencarrier
ExecStart=/usr/local/bin/opencarrier start
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Environment
EnvironmentFile=-/etc/opencarrier/env
WorkingDirectory=/var/lib/opencarrier

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/opencarrier
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=false
RestrictRealtime=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
